# conf/caddy.conf
[program:caddy]
command=caddy -conf "/home/rtorrent/.caddy/Caddyfile"
autostart=true
autorestart=true
user = rtorrent
environment = HOME="/home/rtorrent/", USER="rtorrent"
