# conf/rt.conf
[program:rtorrent]
user=rtorrent
directory=/home/rtorrent
command=rtorrent
redirect_stderr=true
environment=HOME="/home/rtorrent",PWD="/home/rtorrent",LOGNAME="rtorrent",USER="rtorrent",TERM=xterm
